# rootfs/etc/s6-overlay/s6-rc.d/init-vpn/up
with-contenv
/scripts/init-scripts/vpn.sh