# rootfs/etc/s6-overlay/s6-rc.d/init-environment/up
with-contenv
/scripts/init-scripts/environment.sh