# rootfs/etc/s6-overlay/s6-rc.d/init-qbittorrent/up
with-contenv
/scripts/init-scripts/qbittorrent.sh