# rootfs/etc/s6-overlay/s6-rc.d/init-network/up
with-contenv
/scripts/init-scripts/network.sh