# rootfs/etc/s6-overlay/s6-rc.d/service-healthcheck/run
#!/command/with-contenv bash
# shellcheck shell=bash

/scripts/service-scripts/healthcheck.sh